# com.openclaw.openai-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.openai-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/python3</string>
        <string>/path/to/openai-openclaw-bridge-streaming.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/openclaw/workspace</string>
    <key>StandardOutPath</key>
    <string>/path/to/openclaw/workspace/logs/openai-bridge-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/openclaw/workspace/logs/openai-bridge-launchd.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
